# setup/resources/speaker-agent.service
[Unit]
Description=Bluetooth speaker agent

[Service]
ExecStart=PYTHON_DBUS_LOC APP_DIR/speaker-agent.py

[Install]
WantedBy=default.target
